# config/autostart/gnome-terminal.desktop
[Desktop Entry]
Type=Application
Exec=gnome-terminal
Hidden=false
X-GNOME-Autostart-enabled=true
Name=gnome-terminal
